# contrib/systemd/f6-notifier.service
[Unit]
Description=Fractal6 Notifier (MP)
ConditionPathExists=/home/admin/fractal6
After=network.target

[Service]
Type=simple
User=admin
Group=admin

WorkingDirectory=/home/admin/fractal6
ExecStart=/home/admin/fractal6/f6 notifier

Restart=always
RestartSec=5s

StandardOutput=append:/var/log/fractal6/notifier.log
StandardError=append:/var/log/fractal6/err.log

[Install]
WantedBy=multi-user.target
